# deploy/feishu-scheduler.service
[Unit]
Description=Feishu Scheduler (大修排班)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/feishu-scheduler
EnvironmentFile=/opt/feishu-scheduler/.env
ExecStart=/opt/feishu-scheduler/.venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:5000 \
    --access-logfile /var/log/feishu-scheduler/access.log \
    --error-logfile /var/log/feishu-scheduler/error.log \
    app:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
